# smartupdate.conf
[program:smartupdate]
command=/bin/bash -c "smartupdate"
